# meta-jetson-customnux-dashboard/recipes-core/vision-ai-project/vision-backend/vision-backend.service
[Unit]
Description=Vision AI Backend Service
After=network.target

[Service]
ExecStart=/usr/bin/vision_backend
Restart=always
RestartSec=5
LogsDirectory=vision

[Install]
WantedBy=multi-user.target
